# hbos-bluetooth-service-0.1.0/DEBIAN/postinst
#!/bin/sh
set -e
systemctl enable hbos-bluetooth.service
systemctl start hbos-bluetooth.service
