# hbos-bluetooth-service-0.1.0/DEBIAN/prerm
#!/bin/sh
set -e
systemctl stop hbos-bluetooth.service || true
systemctl disable hbos-bluetooth.service || true
